# app/src/main/java/.dailySymlink.cron
30 0 * * * cd ~/public_html/Resources/Android/symlinkScripts; ./.symlinkLatestVersion.sh
