# src/.dailySymlink.cron
30 0 * * * cd ~/public_html/Resources/Android/; ./.symlinkLatestVersion.sh
