# src/.dailySymlink.cron
30 * * * * ~/public_html/Resources/Android/.symlinkLatestVersion.sh
