# src/.dailySymlink.cron
30 * * * * cd ~/public_html/Resources/Android/; ./.symlinkLatestVersion.sh
